# pyget.timer
[Unit]
Description=Run pyget.py every 10 minutes

[Timer]
OnCalendar=*:0/10
Persistent=true

[Install]
WantedBy=timers.target
